# etc/sv/fsrest-jekor.com/run
#!/bin/sh
exec chpst -uhttp:http fsrest /home/jekor/jekor.com/www 10035 2>&1
