# etc/sv/fsrest-jekor.com/run
#!/bin/sh
exec chpst -uwww:www /usr/local/bin/fsrest /home/jekor/jekor.com/www 10035 2>&1
